# playbooks/files/Xvfb.service
[Service]
ExecStart=/usr/bin/Xvfb :99 -ac
