# dotfiles/config/systemd/user/Xvfb.service
[Service]
ExecStart=/usr/bin/Xvfb :99 -ac

[Install]
WantedBy=multi-user.target
